# docker/supervisor/conf.d/openvpn-go-api.conf
[program:openvpn-go-api]
; 容器内 API 必须以 root 运行：它要写 root 属主的 /etc/openvpn、签发证书、并通过 supervisorctl 管理 openvpn-server。
; main() 启动时会跑 InitCore（goose 迁移 + seed 超管 + 客户端同步），依赖的 DB_*/JWT_SECRET/SUPERADMIN_* 从 supervisord 继承的容器环境读取。
command=/app/openvpn-go web --port %(ENV_WEB_PORT)s
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=3
user=root
environment=GIN_MODE="release",OPENVPN_CONFIG_DIR="/etc/openvpn"
redirect_stderr=true
stdout_logfile=/var/log/supervisor/openvpn-go-api.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/openvpn-go-api-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
priority=200
